# steam-order-bot.service
[Unit]
Description=Steam buy-order Telegram notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Замени на своего пользователя и путь на VPS (пути ДОЛЖНЫ быть абсолютными):
User=YOUR_VPS_USER
WorkingDirectory=/home/YOUR_VPS_USER/steam_order_bot
EnvironmentFile=/home/YOUR_VPS_USER/steam_order_bot/.env
ExecStart=/home/YOUR_VPS_USER/steam_order_bot/venv/bin/python bot.py
# Чтобы логи (print) шли в journal сразу, а не пачкой:
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
